# supervisor.conf
[supervisord]
nodaemon=true

[program:clip-server]
command=python -m clip_server /root/worker-template/torch-flow.yml

[program:runpod-handler]
command=python -u /handler.py